# files/etc/uci-defaults/99-auto-resize
#!/bin/sh

LOGFILE="/root/auto-resize.log"
echo "===== EXT4 自动扩容开始 $(date) =====" > $LOGFILE

# 获取根目录设备节点
ROOT_DEV=$(findmnt -n -o SOURCE /)
[ -z "$ROOT_DEV" ] && echo "未找到根设备" >> $LOGFILE && exit 0

# 检查是否为 ext4
FSTYPE=$(blkid -s TYPE -o value $ROOT_DEV)
if [ "$FSTYPE" != "ext4" ]; then
    echo "非 ext4 文件系统，跳过: $FSTYPE" >> $LOGFILE
    exit 0
fi

# 解析磁盘和分区号
case "$ROOT_DEV" in
    *nvme*|*mmcblk*)
        DISK=$(echo $ROOT_DEV | sed 's/p[0-9]*$//')
        PART=$(echo $ROOT_DEV | sed 's/.*p//')
        ;;
    *)
        DISK=$(echo $ROOT_DEV | sed 's/[0-9]*$//')
        PART=$(echo $ROOT_DEV | sed 's/.*[^0-9]//')
        ;;
esac

echo "Root device: $ROOT_DEV" >> $LOGFILE
echo "Disk: $DISK, Partition: $PART" >> $LOGFILE

# 检查剩余空间
DISK_SIZE=$(cat /sys/block/$(basename $DISK)/size)
PART_SIZE=$(cat /sys/class/block/$(basename $ROOT_DEV)/size)

if [ "$DISK_SIZE" -le "$PART_SIZE" ]; then
    echo "没有检测到额外空间，无需扩容。" >> $LOGFILE
    exit 0
fi

echo "检测到额外空间，准备扩容分区..." >> $LOGFILE

# 步骤 1: 扩充分区表
echo "准备扩容分区..." >> $LOGFILE

# 优先尝试 growpart
if command -v growpart >/dev/null 2>&1; then
    growpart "$DISK" "$PART" >> $LOGFILE 2>&1
# 如果没有 growpart，强制使用 parted
elif command -v parted >/dev/null 2>&1; then
    echo "使用 parted 进行扩容..." >> $LOGFILE
    # -s 表示脚本模式，100% 表示将分区扩展到磁盘末尾
    parted -s "$DISK" resizepart "$PART" 100% >> $LOGFILE 2>&1
else
    echo "错误：未找到 growpart 或 parted，无法扩容分区！" >> $LOGFILE
    exit 1
fi

# 步骤 2: 强制内核重新读取分区表
echo "刷新分区表..." >> $LOGFILE
partprobe "$DISK" >> $LOGFILE 2>&1 || partx -u "$ROOT_DEV" >> $LOGFILE 2>&1

# 步骤 3: 扩容文件系统
echo "执行 resize2fs..." >> $LOGFILE
resize2fs "$ROOT_DEV" >> $LOGFILE 2>&1

echo "===== 扩容完成 =====" >> $LOGFILE
df -h / >> $LOGFILE
